# is24-metricd-1.1/.svn/text-base/metric.svn-base
#! /bin/bash
#
### BEGIN INIT INFO
# Provides: testdaemon
# Required-Start: 
# Should-Start: 
# Required-Stop: 
# Should-Stop:
# Default-Start:  3 5
# Default-Stop:   0 1 2 6
# Short-Description: Daemon parsing Nagios/Icinga Perfdata and pushing it to graphite
# Description:    Starts metricd
### END INIT INFO

case "$1" in
  start)
    echo "Starting server"
    # Start the daemon 
    python /data/is24/metric/metricd.py start
    ;;
  stop)
    echo "Stopping server"
    # Stop the daemon
    python /data/is24/metric/metricd.py stop
    ;;
  restart)
    echo "Restarting server"
    python /data/is24/metric/metricd.py restart
    ;;
  *)
    # Refuse to do other stuff
    echo "Usage: /etc/init.d/testdaemon.sh {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
